# 备份文件/deploy/systemd/beingdoing-backend.service
[Unit]
Description=BeingDoing Backend (FastAPI/Uvicorn)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/home/gitclone/BeingDoing/src/backend
EnvironmentFile=/etc/beingdoing.env
ExecStart=/mnt/vdb1/miniconda3/envs/py312/bin/uvicorn app.main:app --host 0.0.0.0 --port 8000
Restart=always
RestartSec=3
TimeoutStartSec=60
TimeoutStopSec=30
KillMode=control-group
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
